# dotfiles-overrides/.config/systemd/user/kanata.service
[Unit]
Description=Kanata keyboard remapper
Documentation=https://github.com/jtroo/kanata
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/kanata -c %h/.config/kanata/config.kbd
Restart=on-failure
RestartSec=1s

[Install]
WantedBy=default.target