# dotfiles-overrides/.config/systemd/user/kanata.service
[Unit]
Description=Kanata keyboard remapper
Documentation=https://github.com/jtroo/kanata
After=graphical-session.target
StartLimitBurst=3
StartLimitIntervalSec=60s

[Service]
Type=simple
Environment=WAYLAND_DISPLAY=wayland-1
Environment=XDG_RUNTIME_DIR=/run/user/1000
# Pre-start check to validate devices are available
ExecStartPre=/bin/bash -c 'source %h/omarchy-pbjorklund/bin/find-keyboard-controller.sh; detect_keyboards_for_kanata'
ExecStart=/usr/bin/kanata -c %h/.config/kanata/config.kbd
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target